# root/usr/lib/systemd/system/p3scan.service.d/nethserver.conf
[Service]
User=root
Group=clamupdate
RestartSec=3
Restart=on-success
